# deploy/vm/isdapresyo-backend.service
[Unit]
Description=IsdaPresyo Backend (Node/Express)
After=network.target

[Service]
Type=simple
User=isdapresyo
WorkingDirectory=/opt/isdapresyo/backend

# Environment variables (edit these values on the server)
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=SERVE_FRONTEND=false
Environment=PGSSL=true
# REQUIRED:
# Environment=DATABASE_URL=postgresql://...
# Environment=JWT_SECRET=...
# Environment=FRONTEND_ORIGIN=https://your-site.netlify.app

ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=3

# Hardening (safe defaults)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
